# debian/postinst
#!/bin/bash
#

systemctl daemon-reload
systemctl unmask galliumos-apollolake
systemctl enable galliumos-apollolake
systemctl start galliumos-apollolake

exit 0
